# services/photo_booth_gui.service
[Unit]
Description=Photo Booth GUI
After=graphical.target

[Service]
User=colin
Environment=DISPLAY=:0
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/games:/usr/games
ExecStart=bash /home/colin/simple_pi_photobooth/launch.sh
Restart=always
RestartSec=1
WorkingDirectory=/home/colin

[Install]
WantedBy=graphical.target
